# install/supervisord.ini
[supervisord]
nodaemon=true
logfile=/logs/supervisord-SOAPservice.log
logfile_maxbytes=50MB
logfile_backups=4
loglevel=info
pidfile=/var/run/supervisord.pid

[program:conf-register]
command=node /conf-register/conf-register
autorestart=unexpected
stdout_logfile=/logs/conf-register.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7

[program:SOAPservice]
directory=/src/
command=node SOAPservice.js -c /config/config.json
stdout_logfile=/logs/SOAPservice.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=7
stderr_logfile=/logs/SOAPservice-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
